# includes/etc/supervisor/conf.d/circusd.conf
[program:circusd]
command = /usr/bin/circusd /etc/circus/conf.d/taiga.ini
stdout_events_enabled=true
stderr_events_enabled=true
autorestart = true
startretries = 10